# backend/docker/supervisor/production.conf
; Production supervisord.conf - runs nginx, php-fpm, and the queue worker in
; one container. Extends the same rationale as the local dev
; docker/supervisor/supervisord.conf: App\Core\Media\Concerns\HasMediaConversions
; queues its thumb/medium conversions, so a worker has to be running
; alongside the web process or PropertyResource/PostResource image URLs
; never resolve.
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# [program:queue-worker]
# ; --max-time forces a periodic restart so a worker never runs long enough
# ; to accumulate a slow memory leak from a badly-behaved job.
# command=php /var/www/html/artisan queue:work --tries=3 --sleep=3 --max-time=3600
# autostart=true
# autorestart=true
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
