# backend/docker/supervisor/supervisord.conf
; Replaces Sail's default /etc/supervisor/conf.d/supervisord.conf (bind-mounted
; over it via compose.yaml) to add a queue worker alongside the web process.
; The container's own start-container script hardcodes
; `supervisord -c /etc/supervisor/conf.d/supervisord.conf`, so a *separate*
; conf.d file relying on the top-level supervisord.conf's [include] glob is
; never actually loaded - this file has to be the one in that exact path.
;
; The queue worker matters because App\Core\Media\Concerns\HasMediaConversions
; registers its thumb/medium WebP conversions as .queued() (PROJECT_RULES.md
; §21: "Queue media conversions... never block the request/response cycle on
; them") - without a worker consuming that queue, PropertyResource's
; `images`/`cover_image` URLs point at conversions that never get generated.
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php]
command=%(ENV_SUPERVISOR_PHP_COMMAND)s
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:queue-worker]
command=/usr/bin/php /var/www/html/artisan queue:work --tries=3 --sleep=3
user=%(ENV_SUPERVISOR_PHP_USER)s
environment=LARAVEL_SAIL="1"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
